# supervisord.conf
[program:gunicorn]
command = gunicorn -c /etc/gunicorn/gunicorn.conf ${DJANGO_APP}.wsgi

[program:nginx]
command = /usr/sbin/nginx