# supervisord.conf
[program:gunicorn]
directory = /usr/django/app
command = gunicorn -c /etc/gunicorn/gunicorn.conf ${DJANGO_APP}.wsgi

[program:nginx]
command = /usr/sbin/nginx